# infra/hostinger/network/systemd/streamx-cloudflare-realip-refresh.timer
[Unit]
Description=Refresh StreamX Cloudflare real IP nginx config daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=streamx-cloudflare-realip-refresh.service

[Install]
WantedBy=timers.target
