# infra/hostinger/systemd/hostinger-backup.timer
[Unit]
Description=Daily StreamX Hostinger backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=hostinger-backup.service

[Install]
WantedBy=timers.target
